# playbook/roles/enforcer/templates/supervisor/app.conf
[program:enforcer]
command=nodejs /opt/saasbox/app/enforcer.js
directory=/opt/saasbox/app
stdout_logfile=/var/log/supervisor/%(program_name)s.log
redirect_stderr=true
autorestart=true
user=root
stopasgroup=true
environment=NODE_ENV=production
